# crontab
* * * * * cd /usr/src/app && python execute_cron.py